# output/crontab.d/localhost.testuser.Landing_Zone.cron
# Update from github landingzones DO NOT manually adjust
# Generated cron file for localhost system, user testuser
# put me in $HOME/crontab.d/Landing_Zone.cron
# Activate cron with:
# `cat $HOME/crontab.d/*.cron | crontab -`
# All .cron files should be found or linked at `$HOME/crontab.d`
SHELL=/bin/sh
PATH=/usr/bin:/bin
# [transfer_001] Transfer from input/* to output/
*/15 * * * * /bin/sh output/scripts/transfer_001.sh
